# scripts/openclaw.service
# =============================================================================
# systemd 服务单元 — OpenClaw Gateway
# =============================================================================
# 安装步骤:
#   1. 编辑此文件，将 <YOUR_USER> 替换为你的 Linux 用户名
#   2. 将 WorkingDirectory 改为你的项目实际路径
#   3. sudo cp scripts/openclaw.service /etc/systemd/system/openclaw.service
#   4. sudo systemctl daemon-reload
#   5. sudo systemctl enable openclaw     # 开机自启
#   6. sudo systemctl start openclaw      # 立即启动
#
# 常用命令:
#   sudo systemctl status openclaw       # 查看状态
#   sudo systemctl restart openclaw      # 重启
#   sudo journalctl -u openclaw -f       # 实时日志
# =============================================================================

[Unit]
Description=OpenClaw Gateway - 自媒体创作助手
Documentation=https://docs.clawd.bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
WorkingDirectory=/home/<YOUR_USER>/creator

# 加载环境变量
EnvironmentFile=/home/<YOUR_USER>/.env.creator

# 启动命令
ExecStart=/usr/bin/openclaw start
ExecReload=/bin/kill -HUP $MAINPID

# 自动重启策略
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# 资源限制
LimitNOFILE=65535
MemoryMax=6G

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

# 安全加固
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/<YOUR_USER>/creator /home/<YOUR_USER>/.openclaw /tmp

[Install]
WantedBy=multi-user.target
